# capsule.service
[Unit]
Description=Capsule service

[Service]
Type=simple
ExecStart=/usr/sbin/capsule

[Install]
WantedBy=default.target
